# pecron-monitor.service
[Unit]
Description=Pecron Battery Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# --- EDIT THESE TO MATCH YOUR SYSTEM ---
# User: the account that owns the pecron-monitor directory
# WorkingDirectory / ExecStart: path where you cloned the repo
User=pi
WorkingDirectory=/home/pi/pecron-monitor
ExecStart=/home/pi/pecron-monitor/run.sh
# ----------------------------------------
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
